# supervisord.conf
[supervisord]
nodaemon=true

[program:scrapyd]
command=scrapyd -d /Work/searchSpiders

[program:scrapydeploy]
command=scrapyd-deploy default -p searchSpiders
